# conf/supervisor/check_already_testing.conf
[program:check_already_testing]
command = /usr/bin/skynet_exo /etc/supervisor/conf.d/check_already_testing.conf
process_name = %(program_name)s_%(process_num)s
numprocs = 1
user = bossmaintainer
umask = 022
autostart = true
autorestart = true
startsecs = 5
startretries = 100
stopwaitsecs = 10
redirect_stderr = true
stdout_logfile = /var/log/supervisor/%(program_name)s_%(process_num)s.log
stderr_logfile = off
environment = PYTHONUNBUFFERED=1,HOME="/home/bossmaintainer",USER="bossmaintainer"

[participant]
name = check_already_testing
queue = check_already_testing
regexp = check_already_testing
code = /usr/share/boss-skynet/check_already_testing.py

