# conf/supervisor/check_has_valid_repo.conf
[program:check_has_valid_repo]
command = /usr/bin/skynet_exo /etc/supervisor/conf.d/check_has_valid_repo.conf
process_name = %(program_name)s_%(process_num)s
numprocs = 1
user = bossmaintainer
umask = 022
autostart = true
autorestart = true
startsecs = 5
startretries = 100
stopwaitsecs = 10
redirect_stderr = true
stdout_logfile = /var/log/supervisor/%(program_name)s_%(process_num)s.log
stderr_logfile = off
environment = PYTHONUNBUFFERED=1,HOME="/home/bossmaintainer",USER="bossmaintainer"

[participant]
name = check_has_valid_repo
queue = check_has_valid_repo
regexp = check_has_valid_repo
code = /usr/share/boss-skynet/check_has_valid_repo.py

