# conf/supervisor/filter_test_packages.conf
[program:filter_test_packages]
command = /usr/bin/skynet_exo /etc/supervisor/conf.d/filter_test_packages.conf
process_name = %(program_name)s_%(process_num)s
numprocs = 1
user = bossmaintainer
umask = 022
autostart = true
autorestart = true
startsecs = 5
startretries = 100
stopwaitsecs = 10
redirect_stderr = true
stdout_logfile = /var/log/supervisor/%(program_name)s_%(process_num)s.log
stderr_logfile = /var/log/supervisor/%(program_name)s_%(process_num)s.stderr.log
environment = PYTHONUNBUFFERED=1,HOME="/home/bossmaintainer",USER="bossmaintainer"

[participant]
name = filter_test_packages
queue = filter_test_packages
regexp = filter_test_packages
code = /usr/share/boss-skynet/filter_test_packages.py

