# conf/supervisor/l10n_update_ts.conf
[program:l10n_update_ts]
command = /usr/bin/skynet_exo /etc/supervisor/conf.d/l10n_update_ts.conf
process_name = %(program_name)s_%(process_num)s
numprocs = 1
user = bossmaintainer
umask = 022
autostart = true
autorestart = true
startsecs = 5
startretries = 100
stopwaitsecs = 10
redirect_stderr = true
stdout_logfile = /var/log/supervisor/%(program_name)s_%(process_num)s.log
stderr_logfile = off
environment = PYTHONUNBUFFERED=1,HOME="/home/bossmaintainer",USER="bossmaintainer"

[participant]
name = l10n_update_ts
queue = l10n_update_ts
regexp = l10n_update_ts
code = /usr/share/boss-skynet/l10n_update_ts.py

