# conf/supervisor/wait_for_repo_published.conf
[program:wait_for_repo_published]
command = /usr/bin/skynet_exo /etc/supervisor/conf.d/wait_for_repo_published.conf
process_name = %(program_name)s_%(process_num)s
numprocs = 1
user = bossmaintainer
umask = 022
autostart = true
autorestart = true
startsecs = 5
startretries = 100
stopwaitsecs = 10
redirect_stderr = true
stdout_logfile = /var/log/supervisor/%(program_name)s_%(process_num)s.log
stderr_logfile = off
environment = PYTHONUNBUFFERED=1,HOME="/home/bossmaintainer",USER="bossmaintainer"

[participant]
name = wait_for_repo_published
queue = wait_for_repo_published
regexp = wait_for_repo_published
code = /usr/share/boss-skynet/wait_for_repo_published.py

